# conf/crond.ini
[program:crond]
command=crond -L /var/log/cron/cron.log -f
