# conf/webmin.ini
[program:webmin]
command=/etc/webmin/start --nofork
